# home-g/.config/systemd/user/keyd-reload.service
[Unit]
Description=Labwc config auto reload

[Service]
ExecStart=/bin/sudo /bin/keyd reload
Type=oneshot

[Install]
WantedBy=default.target
