# home-g/.config/systemd/user/keyd-reload.path
[Unit]
Description=Keyd config auto reload

[Path]
PathChanged=%h/.config/keyd.conf

[Install]
WantedBy=default.target
