# home-g/.config/systemd/user/labwc-reload.path
[Unit]
Description=Labwc config auto reload

[Path]
PathChanged=%h/.config/labwc/menu.yml
PathChanged=%h/.config/labwc/rc.yml
PathChanged=%h/.config/labwc/themerc-override

[Install]
WantedBy=default.target
